# templates/etc/init/node-exporter.conf.erb
# node-exporter - export node information daemon
#

description	"node-exporter"

start on runlevel [2345]
stop on runlevel [016]

respawn

script
    . /etc/default/node-exporter
    /usr/bin/node_exporter -web.listen-address :<%= @exporter_port %>
end script
